# modules/systemd-user/services/caffeine.service
[Unit]
Description=Caffeine

[Service]
ExecStart=/usr/bin/caffeine

[Install]
WantedBy=default.target
