# modules/systemd-user/services/dunst.service
[Unit]
Description=Dunst Notification Daemon

[Service]
ExecStart=/usr/bin/dunst

[Install]
WantedBy=default.target
