# modules/systemd-user/services/pasystray.service
[Unit]
Description=Pasystray

[Service]
ExecStart=/usr/bin/pasystray

[Install]
WantedBy=default.target
