# modules/systemd-user/services/taskwarrior-sync.timer
[Unit]
Description=Periodically sync tasks with taskserver

[Timer]
OnCalendar=*:0/30
