# modules/systemd-user/services/taskwarrior-sync.service
[Unit]
Description=Syncs tasks with taskserver

[Service]
Type=oneshot
ExecStart=/usr/bin/task sync
